# resources/service_sample.conf
start on filesystem

script
    exec /path/to/your/deploy/file
end script
